# luci-app-huawei-hilink/root/etc/uci-defaults/99-huawei-hilink
#!/bin/sh

if [ ! -L /www/cgi-bin/huawei_hilink ]; then
    ln -s /usr/bin/huawei_hilink /www/cgi-bin/huawei_hilink
fi

sleep 1

chmod +x /usr/bin/huawei_hilink
