# luci-app-huawei-hilink/root/etc/uci-defaults/99-huawei-hilink
#!/bin/sh

ln -s /usr/bin/huawei_hilink /www/cgi-bin/huawei_hilink

sleep 1

chmod +x /usr/bin/huawei_hilink
